# depfu/crontab
30 * * * * root cd /home/bot && bundle.ruby2.4 exec ruby.ruby2.4 runner.rb >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
