# pull_request_package/obs-package-bot.timer
[unit]
Description=obs-package-bot timer

[Timer]
OnBootSec=10m
OnUnitInactiveSec=5m

[Install]
WantedBy=basic.target
